# ansible/roles/code-server/templates/code-server.service.j2
[Unit]
Description=code-server - Web-based VS Code IDE
Documentation=https://github.com/coder/code-server
After=network.target

[Service]
Type=exec
User={{ roughneck_user }}
Group={{ roughneck_user }}
WorkingDirectory=/home/{{ roughneck_user }}
ExecStart=/usr/bin/code-server --bind-addr 0.0.0.0:10000 --auth password --config /home/{{ roughneck_user }}/.config/code-server/config.yaml /home/{{ roughneck_user }}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
